# var/conf-files/crontab
SHELL=/bin/bash

# m  h   dom mon dow command
 12  *   *   *   *   . @@WORKDIR_DEV@@/bash-env; cron-hourly
